# Saily.Daemon/Package/DEBIAN/postinst
# Install permission
chown root:wheel /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
chmod 0644 /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
chmod +x /usr/bin/Saily.Daemon

# Load Daemon
launchctl unload /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist
launchctl load /Library/LaunchDaemons/com.lakr233.Saily-Daemon.plist

# Permission
chmod 0755 /Applications/Saily.app/Saily
chown -R root:wheel /Applications/Saily.app
chown -R root:wheel /Applications/Saily.app/*

# Uicache
# uicache /Applications/Saily.app
# Broken uicache on Chimera, sorry to say.
uicache -a
# su mobile -c uicache

# Swap our daemon
killall Saily.Daemon

# Force dpkg to continue
exit 0
